# launchagents/com.user.faktory-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.faktory-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/Projects/union/workers && bundle exec faktory-worker -q work-web -q any -r ./worker.rb</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Ensure Ruby/bundle are found -->
        <key>PATH</key>
        <string>/Users/jeff.casimir/.local/share/mise/shims:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Run at load and keep running -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Restart if it crashes -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/faktory-worker.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/faktory-worker.error.log</string>
</dict>
</plist>
